# build/src/supervisord.conf
[supervisord]
nodaemon=true

[program:nodejs]
command=node client.js
stderr_logfile = %(here)s/nodejs-stderr.log
stdout_logfile = %(here)s/nodejs-stdout.log

[program:xl2tpd]
command=/usr/src/app/init.sh
stderr_logfile = %(here)s/xl2tpd-stderr.log
stdout_logfile = %(here)s/xl2tpd-stdout.log