# build/src/supervisord.conf
[supervisord]
nodaemon=true

[program:VPN_userManagment]
command=node src/client.js
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0


[program:VPN_openvpn]
command=/usr/src/app/openvpn.sh
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
